# scripts/check_mk.cron
*/5 * * * * root rm -f /etc/check_mk/oracle_SID_.cache
