# scripts/check_mk.cron
*/5 * * * * root rm -f "/etc/check_mk/oracle_{{ item.0.oracle_db_name }}.cache"
